# rules/sigma/persistence/scheduled-task-run-key-service.yml
title: Persistence via Scheduled Task, Run Key, or Service
id: det-soc-018
status: production
description: |
  Detects scheduled task creation, registry Run key modification,
  or service creation with suspicious command/path indicators.
  Red Canary 2025: scheduled tasks are the most prevalent
  persistence technique.
author: Ridgeline Cyber Detection Engineering
date: 2026/05/17
tags:
  - attack.persistence
  - attack.t1053.005
  - attack.t1547.001
  - attack.t1543.003
logsource:
  category: process_creation
  product: windows
detection:
  schtasks:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
    CommandLine|contains:
      - '\Temp\'
      - '\AppData\'
      - 'powershell'
      - 'rundll32'
  condition: schtasks
falsepositives:
  - Software deployment tools creating scheduled tasks
  - Group Policy applying scheduled tasks at scale
  - Legitimate installer creating Run key entries
level: high
